# netmaker/root/etc/hotplug.d/iface/99-netclient
#!/bin/sh
[ "$ACTION" = ifup ] || exit 0
sleep 60
netok=$(ping 192.168.88.254 -c3 |grep "time=")
if [ "$netok" ]; then
   return
else
   netclient daemon &
fi
